# scripts/mqtt-broker.service
[Unit]
Description=Simple MQTT Broker
After=network.target

[Service]
Type=simple
User=mqtt
Group=mqtt
Environment=NODE_ENV=production
Environment=MQTT_PORT=1883
Environment=WEBSOCKET_PORT=8080
Environment=MQTT_HOST=0.0.0.0
WorkingDirectory=/opt/simple-mqtt
ExecStart=/usr/bin/node dist/broker/example.js
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
